# ansible/roles/setup-server-services/templates/cider-ci_pre-configure.service
[Unit]
Description=Cider-CI Pre-Configure

[Service]
Type=oneshot
ExecStart=/cider-ci/bin/pre-configure
RemainAfterExit=yes
