# demo/test.service
[Unit]
Description=Test

[Service]
Type=oneshot
ExecStart=/cider-ci/docker/test.sh
RemainAfterExit=yes
